# deploy/rhyme-backend.service
[Unit]
Description=RHYME Music Server
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/rhyme/backend
EnvironmentFile=/opt/rhyme/backend/.env
ExecStart=/opt/rhyme/backend/venv/bin/uvicorn app.main:app --host 127.0.0.1 --port 8001
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
